# upstart/docker-influx.conf
description "InfluxDb Docker container"
author "Ian Scrivener - 7/8/2016"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a influxdb
end script